# build/rootfs/etc/profile.d/03-aliases.sh
alias code="codium"
alias k="kubectl"
alias vi="vim"
